# build-docker-images/wazuh-manager/config/supervisord.conf
[supervisord]
user=root
logfile=/supervisord.log
pidfile=/var/run/supervisord.pid
nodaemon=true

[program:wazuh]
command=bash -c "./var/ossec/bin/*-control start"
